# relysyspackages/snmppackages/files/etc/uci-defaults/01_snmp_firstboot
#!/bin/sh

. /lib/functions.sh
sleep 5

board_name=$(cat /tmp/sysinfo/board_name)
                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                             
if echo "$board_name" | grep -qE "(Silbo_RC44|Silbo_RB44|Silbo_RD44-3|Silbo_RD44-A|Silbo_GRD44-3|Silbo_GRD44B-5R|Silbo_RF44|Silbo_IA44-A|Silbo_IA44-B|Silbo_IA44-C|Silbo_IAC44-A)";
then
	#DI=0,AI=0,Temp=0,DIOTrap=0
	cp /snmp_boardinfo/Silbo_0_0_0_0/snmp_boardinfo.txt /etc/snmp/

elif echo "$board_name" | grep -qE "(Silbo_IAC44)";
then
	#DI=2,AI=0,Temp=0,DIOTrap=2
	cp /snmp_boardinfo/Silbo_DIO_0_0_DIOtrap/snmp_boardinfo.txt /etc/snmp/
	#cp /snmp_boardinfo/usr/share/rpcd/menu.d/snmpmenu.json /usr/share/rpcd/menu.d/

elif echo "$board_name" | grep -qE "(Silbo_IAB44-B|Silbo_IAB44-C)";
then
	#DI=4,AI=2,Temp=0,DIOTrap=4
	cp /snmp_boardinfo/Silbo_DIO_AI_0_DIOtrap/snmp_boardinfo.txt /etc/snmp/
	#cp /snmp_boardinfo/usr/share/rpcd/menu.d/snmpmenu.json /usr/share/rpcd/menu.d/
	
elif echo "$board_name" | grep -qE "(Silbo_GD44-A|Silbo_GD44-B|Silbo_GD44-C)";
then
	#DI=4,AI=0,Temp=8,DIOTrap=4
	cp /snmp_boardinfo/Silbo_DIO_0_Temp_DIOtrap/snmp_boardinfo.txt /etc/snmp/
	#cp /snmp_boardinfo/usr/share/rpcd/menu.d/snmpmenu.json /usr/share/rpcd/menu.d/

elif echo "$board_name" | grep -qE "(Silbo_GD44-D|Silbo_GD44-E)";
then
	#DI=4,AI=2,Temp=8,DIOTrap=4
	cp /snmp_boardinfo/Silbo_DIO_AI_Temp_DIOtrap/snmp_boardinfo.txt /etc/snmp/
	#cp /snmp_boardinfo/usr/share/rpcd/menu.d/snmpmenu.json /usr/share/rpcd/menu.d/
fi

exit 0
